# root/etc/uci-defaults/zz-prometheus-node-exporter-bind
#!/bin/sh

uci set prometheus-node-exporter.main.listen_address '::'
uci commit prometheus-node-exporter
