# packages/meshrc/files/web/etc/uci-defaults/restart_rpcd.sh
#!/bin/sh

# make rpcd aware of the added meshrc commands
/etc/init.d/rpcd restart
